# docs-svc/docs-rsync.timer
[Timer]
OnCalendar=daily

[Install]
WantedBy=default.target
