# docs-svc/docs.service.in
[Unit]
Description=Hosts developer documentation at localhost:10000

[Service]
Type=forking
ExecStartPre=@MKDIR@ -p /tmp/jade/nginx
ExecStart=@NGINX@ -p .dotfiles/docs-svc -e stderr -c nginx.conf
WorkingDirectory=~
NoNewPrivileges=true

[Install]
WantedBy=default.target
