# docs-svc/docs-rsync-pull.service
[Unit]
Description=rsync dev docs from abyss

[Service]
Type=simple
ExecStart=%h/.dotfiles/docs-svc/pull.sh
WorkingDirectory=%h/.dotfiles/docs-svc

