# thinkpad/trackpoint-setup.service
[Unit]
Description=Set up a TrackPoint device

[Service]
Type=oneshot
ExecStart=/usr/local/bin/trackpoint-setup

[Install]
WantedBy=graphical.target
